# config/autostart/my-startup.desktop
[Desktop Entry]
Type=Application
Exec=/home/bluekitty/.local/bin/my-startup.sh
Name=My Startup Apps
Hidden=false
